# deploy/files/supervisor.conf
[program:vodka]
command = gunicorn app:app -b localhost:8000
directory = /vagrant/